# build/nginx-fpm/crontab
* 0 * * * /www/app/bin/console app:pdf
* 0 * * * /www/app/bin/console ajax_uploader:clean 1 -q